# crontab

* * * * * /home/tomh/bin/iphonepics.php >/dev/null 2>&1
